# systemd/wazuh-immutable-store-verify.service
[Unit]
Description=Wazuh Immutable Store - Integrity Verification Service
Documentation=https://github.com/your-repo/wazuh-immutable-store
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wazuh-immutable-store verify
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wazuh-immutable-store-verify

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

# Allow access to necessary directories
ReadOnlyPaths=/tmp/wazuh-archive
ReadOnlyPaths=/mnt/qnap-wazuh
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
